# deploy/pickleball-logs.service
[Unit]
Description=Nightly referee-log backfill (see scraper/harvest_logs.py)
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
# Optional Supabase creds for the warehouse upsert (SUPABASE_URL,
# SUPABASE_SERVICE_KEY). The leading '-' makes it non-fatal if absent, so
# the harvest still runs without them; upload_supabase.py no-ops.
EnvironmentFile=-@REPO@/deploy/pickleball.env
ExecStart=@REPO@/deploy/run_logs_backfill.sh
# No Restart: the harvest is resumable and re-fires tomorrow night anyway.
